# deploy/mneme.service
[Unit]
Description=Mneme Knowledge Base Server
Documentation=https://github.com/agnostos/mneme
After=network.target daimon.service
Wants=daimon.service

[Service]
Type=simple
ExecStart=/usr/local/bin/mneme-api
Environment=MNEME_VAULT_DIR=%h/mneme
Environment=MNEME_BIND=127.0.0.1:3838
Environment=DAIMON_URL=http://127.0.0.1:8090
Environment=RUST_LOG=info
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

# Sandboxing
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/mneme
PrivateTmp=true
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictSUIDSGID=true

[Install]
WantedBy=default.target
